# server/quiz_server.socket
[Unit]
Description=Django quiz gunicorn socket

[Socket]
ListenStream=/run/quiz_server.sock

[Install]
WantedBy=sockets.target
